# build/docker/runit/apisix/finish
#!/bin/bash
# runit finish script for apisix
