# build/docker/runit/isrvd/run
#!/bin/bash
exec 2>&1

# 等待 APISIX Admin API 就绪（etcd 已由 apisix 等待）
echo "[isrvd] Waiting for APISIX..."
for i in $(seq 1 60); do
    if curl -sf --connect-timeout 5 --max-time 10 http://127.0.0.1:9180/apisix/admin/routes &>/dev/null; then
        echo "[isrvd] APISIX is ready"
        break
    fi
    if [ "$i" -eq 60 ]; then
        echo "[isrvd] APISIX not ready after 60s, continuing anyway..."
        break
    fi
    sleep 1
done

echo "[isrvd] Starting isrvd..."
cd /data

echo "[isrvd] Using CONFIG_PATH: $CONFIG_PATH"
exec /usr/local/bin/isrvd
